# 23_icbc/supervisord/codeserver.ini
[program:codeserver]
command=code-server
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/codeserver.err.log
stdout_logfile=/var/log/codeserver.log