# 23_icbc/supervisord/crewai.ini
[program:crewai]
command=ollama serve
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/crewai.err.log
stdout_logfile=/var/log/crewai.log