# 23_icbc/supervisord/flower.ini
[program:flower]
directory=/share/mlops
command=celery -A mlops.app.app.celery_app flower --port=5432
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/flower.err.log
stdout_logfile=/var/log/flower.log
