# 23_icbc/supervisord/mlflow.ini
[program:mlflow]
command=mlflow server --host 0.0.0.0 --port 8082 --backend-store-uri file:/share/mlflow/mlruns --default-artifact-root file:/share/mlflow/artifacts
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/mlflow.err.log
stdout_logfile=/var/log/mlflow.log