# 23_icbc/supervisord/mlops.ini
[program:mlops]
directory=/share/mlops
command=gunicorn --workers 1 --timeout 120 'mlops.run'
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/mlops.err.log
stdout_logfile=/var/log/mlops.log
