# 23_icbc/supervisord/mlopsworker.ini
[program:mlopsworker]
directory=/share/mlops
command=celery -A mlops.app.app.celery_app worker -l info
autostart=true
autorestart=true
startretries=1
stderr_logfile=/var/log/mlopsworker.err.log
stdout_logfile=/var/log/mlopsworker.log