# 23_superset/04_dev/jupyter.ini
[program:jupyter]
command=jupyter lab --allow-root  --ServerApp.tornado_settings='{"headers":{"Content-Security-Policy":"frame-ancestors self *; report-uri /api/security/csp-report"}}' --allow-root --no-browser --ip=0.0.0.0 --port=8888 --ServerApp.token='' --ServerApp.password='' --ServerApp.quit_button=False
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/jupyter.err.log
stdout_logfile=/var/log/jupyter.log