# 18_ubuntu/03_python/django.ini
[program:django]
command=/usr/bin/python3 /share/phoenix/aihf/manage.py runserver 0.0.0.0:8003
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/django.err.log
stdout_logfile=/var/log/django.out.log